# music_player_files/rave.timer
[Unit]
Description=Plays crab rave at 1000 and 2200

[Timer]
OnCalendar=*-*-* 10:00:00
OnCalendar=*-*-* 22:00:00

[Install]
WantedBy=timers.target