# music_player_files/rave.timer
[Unit]
Description=Plays crab rave every 6 hours

[Timer]
Persistent=true
OnBootSec=1
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target