# music_player_files/rave.service
Description=Play Crab Rave

[Service]
ExecStart=sh /home/USER/Music/rave.sh

[Install]
WantedBy=multi-user.target